# etc/crontab/stage_s2-wrk-production.cron.twig
# Managed by deployer
# Cronjobs should go to worker nodes only. For S2 it's s2-wrk-production environment.

# m h dom mon dow command

# run every day at 10:25
25 10 * * * /usr/bin/php{{ php_version }} {{ project_dir }}/bin/console app:cron:demo-lock >/dev/null

# runs every day between 2:00 and 3:00 o'clock (start time (minute) is randomly generated for each node)
{{ random(59) }} 2 * * * /usr/bin/php{{ php_version }} {{ project_dir }}/bin/console app:cron:demo-lock >/dev/null

# runs every 15 minutes between 3:00 and 4:00 o'clock (start time (minute) is randomly generated for each node)
{{ random(14) }}-59/15 3 * * * /usr/bin/php{{ php_version }} {{ project_dir }}/bin/console app:cron:demo-lock >/dev/null
